# network-monitor.service
[Unit]
Description=Network Monitor
After=network.target

[Service]
Type=simple
User=kerv
WorkingDirectory=/path/to/network-monitor
ExecStart=/usr/bin/python3 /path/to/network-monitor/network_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
